# observability-and-management/assets/oci-log-analytics-detections/rules/linux/linux_suspicious_cron_content.yaml
title: Linux Suspicious Cron Job Content
id: f1e2d3c4-5009-4000-8000-000000000009
status: stable
description: Detects cron job entries containing suspicious patterns such as downloading
  payloads, piping to shell, or encoding commands typically used for persistence.
logsource:
  product: linux
  service: audit
detection:
  sel_cron_dir:
    message|contains:
    - /etc/cron
    - /var/spool/cron
  sel_suspicious:
    message|contains:
    - curl|sh
    - wget|sh
    - curl|bash
    - wget|bash
    - base64 -d
  condition: sel_cron_dir and sel_suspicious
level: high
tags:
- attack.persistence
- attack.t1053.003
falsepositives:
- Legitimate automation scripts using curl piped to shell
version: 1
