# observability-and-management/assets/oci-log-analytics-detections/rules/linux/linux_crontab_persistence.yaml
title: Linux Crontab Modification
id: 07b8c9d0-bcde-f012-3456-789012345067
status: stable
description: Detects crontab command execution, which may indicate persistence via scheduled tasks.
logsource:
  product: linux
  service: audit
detection:
  selection:
    process_name: crontab
  suspicious_args:
    command_line|contains:
    - -e
    - -r
    - /tmp/
    - /var/tmp/
    - /dev/shm/
  condition: selection and suspicious_args
level: medium
tags:
- attack.persistence
- attack.t1053.003
version: 1
falsepositives:
- Legitimate scheduled maintenance tasks
- Automated system updates
