# observability-and-management/assets/oci-log-analytics-detections/queries/impacket_atexec_remote_scheduled_task_execution.json
{
  "title": "Impacket atexec Remote Scheduled Task Execution",
  "description": "Detects impacket atexec, which creates a transient scheduled task on a remote host (via the ATSvc/Task Scheduler) to run a command and redirect output to an admin share. Indicators: a task launched by the Task Scheduler running cmd with output redirected to ADMIN$/a temp marker file.\n",
  "query": "('Log Source' = 'Windows Sysmon Events' or 'Log Source' = 'Windows Sysmon Operational Logs' or 'Log Source' = 'SOC Windows Sysmon Logs') and ((('Parent Process Name' like '*svchost.exe' or 'Parent Process Name' like '*taskeng.exe')) and (('Command Line' like '*cmd.exe*' and 'Command Line' like '*2>&1*')) and (('Command Line' like '*ADMIN$*' or 'Command Line' like '*Temp*' or 'Command Line' like '*2>&1*')))",
  "sigma_id": "407543f0-898a-4b2c-bfb3-804bb2a0ede6",
  "level": "high",
  "stig_category": "CAT I",
  "tags": [
    "attack.execution",
    "attack.t1053.005",
    "attack.lateral_movement"
  ],
  "mitre_attack": {
    "tactics": [
      "execution",
      "lateral_movement"
    ],
    "techniques": [
      "T1053.005"
    ]
  },
  "logsource": {
    "product": "windows",
    "service": "",
    "candidates": [
      "Windows Sysmon Events",
      "Windows Sysmon Operational Logs",
      "SOC Windows Sysmon Logs"
    ]
  },
  "falsepositives": [
    "Legitimate scheduled-task based administration"
  ]
}