# app-dev/devops/ansible-jenkins/jenkins-playbook/templates/jenkins.service.j2
[Service]
Environment="JAVA_OPTS=-Djava.awt.headless=true -Dcasc.jenkins.config={{ jenkins_casc_path }} -Djenkins.install.runSetupWizard=false"
Environment="JENKINS_PORT={{ jenkins_port }}"